# app/product/crontab
#crontab
* * * * * root /usr/local/bin/python3 /usr/src/app/daemon.py -d 3 -p 10 -allp 1 >> /var/log/cron.log 2>&1

